# zero-token-cron/launchd/com.dayton.healthcheck.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dayton.healthcheck</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/daytrons/.openclaw/workspace/zero-token-cron/jobs/health_check_job.sh</string>
    </array>
    
    <key>StartInterval</key>
    <integer>21600</integer>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/daytrons/.openclaw/workspace/zero-token-cron/logs/healthcheck.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/daytrons/.openclaw/workspace/zero-token-cron/logs/healthcheck.error.log</string>
    
    <key>WorkingDirectory</key>
    <string>/Users/daytrons/.openclaw/workspace</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
